# 1.2.1/usr/etc/supervisord.conf
[supervisord]

[program:riak]
command=su riak - /usr/sbin/riak console
